# systemd/blog.service
[Unit]
Description=Next.js Blog Application
After=network.target

[Service]
User=your_username
WorkingDirectory=/Users/sawyer/Library/CloudStorage/GoogleDrive-stephenvsawyer@gmail.com/My Drive/github/t3-shipyard/apps/blog
ExecStart=/usr/bin/npm start
Restart=always

[Install]
WantedBy=multi-user.target
