# systemd/codecaddy.service
[Unit]
Description=CodeCaddy Next.js Application
After=network.target

[Service]
User=%i
WorkingDirectory=/home/%i/t3-shipyard/apps/codecaddy
ExecStart=/usr/bin/pm2 start npm --name "codecaddy-app" -- start
ExecStop=/usr/bin/pm2 stop codecaddy-app
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=codecaddy-app

[Install]
WantedBy=multi-user.target
