# systemd/dashboard.service
[Unit]
Description=Next.js Dashboard Application
After=network.target

[Service]
# Replace with your actual user and group
User=www-data
Group=www-data

# Replace with the actual path to your project
WorkingDirectory=/srv/t3-shipyard/apps/dashboard

ExecStart=/usr/bin/pm2 start npm --name "dashboard-app" -- start
Restart=always

[Install]
WantedBy=multi-user.target