# systemd/dashboard.service
[Unit]
Description=Next.js Dashboard Application
After=network.target

[Service]
User=your_username
WorkingDirectory=/Users/sawyer/Library/CloudStorage/GoogleDrive-stephenvsawyer@gmail.com/My Drive/github/t3-shipyard/apps/dashboard
ExecStart=/usr/bin/pm2 start npm --name "dashboard-app" -- start
Restart=always

[Install]
WantedBy=multi-user.target
