# systemd/dunamismax.com.service
[Unit]
Description=Next.js Dunamismax.com Application
After=network.target

[Service]
# Replace with your actual user and group
User=www-data
Group=www-data

# Replace with the actual path to your project
WorkingDirectory=/srv/t3-shipyard/apps/dunamismax.com

ExecStart=/usr/bin/pm2 start npm --name "dunamismax-app" -- start
Restart=always

[Install]
WantedBy=multi-user.target