# systemd/trpchat.service
[Unit]
Description=PM2 process for trpchat application
After=network.target

[Service]
# Replace with your actual user and group
User=www-data
Group=www-data

# Replace with the actual path to your project
WorkingDirectory=/srv/t3-shipyard/apps/trpchat

ExecStart=/usr/bin/pm2 start ecosystem.config.js --env production
ExecReload=/usr/bin/pm2 reload all
ExecStop=/usr/bin/pm2 stop all
Restart=always
RestartSec=3
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=trpchat

[Install]
WantedBy=multi-user.target