# systemd/trpchat.service
[Unit]
Description=PM2 process for trpchat application
After=network.target

[Service]
User=your_username
Group=your_username
WorkingDirectory=/Users/sawyer/Library/CloudStorage/GoogleDrive-stephenvsawyer@gmail.com/My Drive/github/t3-shipyard/apps/trpchat
ExecStart=/usr/bin/pm2 start ecosystem.config.js --env production
ExecReload=/usr/bin/pm2 reload all
ExecStop=/usr/bin/pm2 stop all
Restart=always
RestartSec=3
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=trpchat

[Install]
WantedBy=multi-user.target
